# ansible/roles/imagemaker/templates/imagemaker.service.j2
[Unit]
Description=openBalena Image Maker
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User={{ service_user }}
Group={{ service_user }}
WorkingDirectory={{ install_root }}/imagemaker
Environment=SERVICE_HOME=/home/{{ service_user }}
Environment=DNS_TLD={{ dns_tld }}
Environment=PUBLIC_TLD={{ public_tld }}
Environment=BALENARC_BALENA_URL={{ public_tld }}
Environment=OPENBALENA_ROOT_CA=/usr/local/share/ca-certificates/openbalena-root-ca.crt
Environment=IMAGEMAKER_DIR={{ install_root }}/imagemaker
# Shared secret required to use the service (images embed provisioning creds). Open the UI as
# http://127.0.0.1:8090/?token=<value>. Leave blank only for a pure loopback/tunnel-only setup.
Environment=IMAGEMAKER_TOKEN={{ imagemaker_token | default('') }}
ExecStart=/usr/bin/node {{ install_root }}/imagemaker/server.js
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
